# telegram.ini
[program:myapp]
command=/data/telegram_python/ll_env/bin/python3 /data/telegram_python/main.py
directory=/data/telegram_python
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/telegram/myapp.log
stderr_logfile=/var/log/telegram/myapp_error.log